# install/clipd.service
[Unit]
Description=LinVClipBoard Clipboard Daemon
Documentation=https://github.com/DreamerX00/LinVClipBoard
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/clipd
Restart=always
RestartSec=3
Environment=RUST_LOG=info

# Graceful shutdown: SIGTERM then SIGKILL after 5s
TimeoutStopSec=5
KillSignal=SIGTERM

# Resource limits
MemoryMax=50M
CPUQuota=5%

# Sandboxing (systemd hardening)
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=%h/.local/share/linvclip %h/.config/linvclip %t/linvclip
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
