# install/linvclip-update-check.timer
[Unit]
Description=Check for LinVClipBoard updates weekly

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
