# packaging/prerm
#!/bin/sh
set -e

# Runs as root during dpkg — iterate real user sessions (#6).
echo "Stopping clipd service for active users..."
for user_run in /run/user/[0-9]*; do
    uid=$(basename "$user_run")
    [ "$uid" -ge 1000 ] 2>/dev/null || continue
    username=$(id -nu "$uid" 2>/dev/null) || continue
    su "$username" -c "XDG_RUNTIME_DIR=$user_run DBUS_SESSION_BUS_ADDRESS=unix:path=$user_run/bus systemctl --user stop clipd.service" 2>/dev/null || true
    su "$username" -c "XDG_RUNTIME_DIR=$user_run DBUS_SESSION_BUS_ADDRESS=unix:path=$user_run/bus systemctl --user stop linvclip-update-check.timer" 2>/dev/null || true
    # Only disable services on full remove, not on upgrade
    if [ "$1" = "remove" ]; then
        su "$username" -c "XDG_RUNTIME_DIR=$user_run DBUS_SESSION_BUS_ADDRESS=unix:path=$user_run/bus systemctl --user disable clipd.service" 2>/dev/null || true
        su "$username" -c "XDG_RUNTIME_DIR=$user_run DBUS_SESSION_BUS_ADDRESS=unix:path=$user_run/bus systemctl --user disable linvclip-update-check.timer" 2>/dev/null || true
    fi
done

# Only kill processes on full remove, not on upgrade
if [ "$1" = "remove" ]; then
    pkill -x clipd 2>/dev/null || true
    pkill -x linvclip-ui 2>/dev/null || true
fi
